# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/nginx.err.log
stdout_logfile=/var/log/nginx.out.log
user=root

[program:backend]
command=python run_with_scheduler.py
directory=/app/backend
autostart=true
autorestart=true
stderr_logfile=/var/log/backend.err.log
stdout_logfile=/var/log/backend.out.log
user=root
environment=FLASK_ENV="production",FLASK_APP="app",DATABASE_PATH="/app/backend/instance/finance.db",TINYDB_CONFIG_PATH="/app/backend/instance/config.json",TZ="Asia/Shanghai"

